# debian/upaas-router.upstart
description "uWSGI FastRouter for uPaaS"

start on runlevel [2345]
stop on runlevel [06]

limit nofile 100000 100000

exec uwsgi --ini /etc/upaas/upaas-router.ini
