# debian/upaas-router.upstart
description "uWSGI FastRouter for uPaaS"

start on runlevel [2345]
stop on runlevel [06]

pre-start script
    mkdir -p /run/upaas-router
    chown www-data:www-data /run/upaas-router
end script

exec uwsgi --ini /etc/upaas-router.ini
